# example/root/etc/services.d/example/run
#!/usr/bin/with-contenv sh
node /app/app.js
